# deploy/systemd/frame-core.service
[Unit]
Description=Picture Frame core service (HTTP API, CDP, scheduler, updater)
After=network.target
Wants=network.target

[Service]
Type=simple
User=frame
Group=frame
WorkingDirectory=/opt/frame/current
Environment=NODE_ENV=production
Environment=FRAME_CONFIG=/etc/frame/frame.yaml
Environment=FRAME_STATE_DIR=/opt/frame/state
Environment=FRAME_RELEASES_DIR=/opt/frame/releases
Environment=FRAME_RUNTIME_DIR=/run/frame
# Cage's chromium is launched with --remote-debugging-port=9222 on
# loopback; frame-core connects to it via WsTransport for URL-screen
# tab management (see deploy/launch-chromium.sh + core/src/cdp/).
Environment=FRAME_CDP_HOST=127.0.0.1
Environment=FRAME_CDP_PORT=9222
ExecStart=/usr/bin/node /opt/frame/current/core/dist/index.js
Restart=on-failure
RestartSec=3
LimitNOFILE=8192

# Create /run/frame on each start (tmpfs clears it on reboot); systemd makes
# this before the namespace is set up, so it's safe to reference below.
RuntimeDirectory=frame
RuntimeDirectoryMode=0755

# Hardening (the Wayland session lives in frame-kiosk; core only needs files + net).
# NoNewPrivileges is intentionally OFF: frame-core invokes sudo for the narrow
# allowlist in /etc/sudoers.d/frame (systemctl restart frame-core/kiosk and
# systemctl reboot), and the terminal-in-browser feature spawns a bash that
# also needs to exec those sudo commands. NoNewPrivileges blocks all setuid
# (sudo can't elevate), which broke reboot from both the web UI and terminal.
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/opt/frame /etc/frame /var/log
PrivateTmp=true

[Install]
WantedBy=multi-user.target
