# deploy/systemd/frame-kiosk.service
[Unit]
Description=Picture Frame kiosk (Cage + Chromium driven by frame-core)
After=frame-core.service systemd-user-sessions.service getty@tty1.service
Requires=frame-core.service
Conflicts=getty@tty1.service
# Give up after 5 failures in 120s. Without this, cage failing (e.g. no DRM
# in a VM) traps tty1 in an endless restart loop and the system never
# presents a usable console.
StartLimitIntervalSec=120
StartLimitBurst=5
OnFailure=getty@tty1.service

[Service]
Type=simple
User=frame
Group=frame
PAMName=login
TTYPath=/dev/tty1
TTYVTDisallocate=yes
StandardInput=tty
StandardOutput=journal
StandardError=journal
WorkingDirectory=/home/frame
Environment=WLR_LIBINPUT_NO_DEVICES=1
Environment=XCURSOR_THEME=frame-transparent
Environment=XCURSOR_SIZE=1
ExecStartPre=/opt/frame/current/deploy/wait-for-core.sh
ExecStart=/usr/bin/cage -- /opt/frame/current/deploy/launch-chromium.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
